# .crontab
0 * * * * $HOME/scripts/auto_commit.sh $HOME/notes
